# smb_dir_backup.service
[Unit]
Description=Backup Samba directory


[Service]
Type=simple
ExecStart=/home/serg/scripts/backup


[Install]
WantedBy=multi-user.target